# base/.config/sh/profile.d/50-ccache.sh
export CCACHE_DIR="${XDG_CACHE_HOME:-"$HOME/.cache"}/ccache"
export CCACHE_CONFIGPATH="${XDG_CONFIG_HOME:-"$HOME/.config"}/ccache.config"
