# base/.config/sh/profile.d/50-go.sh
export GOPATH="${XDG_CACHE_HOME:-"$HOME/.cache"}/go"
export GOBIN="${XDG_BIN_HOME:-"$HOME/.local/bin"}"
