# base/.config/sh/profile.d/50-openssl.sh
export RANDFILE="${XDG_CACHE_HOME:-"$HOME/.cache"}/rand"
