# base/.config/sh/profile.d/50-python.sh
export PYTHON_EGG_CACHE="${XDG_CACHE_HOME:-"$HOME/.cache"}/python-eggs"
export PYLINTHOME="${XDG_CACHE_HOME:-"$HOME/.cache"}/pylint"
export PYTHONSTARTUP="${XDG_CONFIG_HOME:-"$HOME/.config"}/python/cli.py"
