# systemd/user/arch-audit.timer
[Unit]
Description=Hourly check for Arch Linux vulnerabilities

[Timer]
OnStartupSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
